# heartbeat/init.mysql-heartbeat.conf
start on starting mysql
stop on stopping mysql

exec mk-heartbeat --update --database=info --user=heartbeat --password=heartsmart

